# librescoot-update-installer-dbc.service
[Unit]
Description=Librescoot Update Installer (DBC)
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
ExecStart=/usr/bin/update-installer \
    --redis-addr=192.168.7.1:6379 \
    --update-key=update:install:dbc \
    --failure-key=update:failure \
    --component=dbc \
    --update-type=blocking \
    --fetch-queue=fetch-requests
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target