# librescoot-update-mdb.service
[Unit]
Description=LibreScoot Update Service (MDB)
After=redis.service
After=network-online.target
After=librescoot-modem.service
After=librescoot-vehicle.service
After=librescoot-version.service
After=librescoot-pm.service
Requires=redis.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/update-service --component=mdb --channel=nightly
Restart=always
RestartSec=10
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=update-service-mdb

# Create data directory if it doesn't exist
ExecStartPre=/bin/mkdir -p /data/ota/mdb

[Install]
WantedBy=multi-user.target