# librescoot-update-mdb.service
[Unit]
Description=LibreScoot Update Service (MDB)
After=redis.service network.target
Requires=redis.service

[Service]
Type=simple
ExecStart=/usr/bin/update-service --component=mdb --channel=nightly
Restart=always
RestartSec=10
User=root
Group=root

# Environment variables
Environment=COMPONENT=mdb
Environment=CHANNEL=nightly

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=update-service-mdb

# Working directory
WorkingDirectory=/data/ota

# Create data directory if it doesn't exist
ExecStartPre=/bin/mkdir -p /data/ota/mdb

[Install]
WantedBy=multi-user.target