# librescoot-update-fetcher.service
[Unit]
Description=Librescoot Update Fetcher Service
After=network-online.target redis.service
Wants=network-online.target redis.service

[Service]
Type=simple
ExecStart=/usr/bin/update-fetcher \
    --redis-addr=192.168.7.1:6379 \
    --download-dir=/data/ota \
    --queue=fetch-requests
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target